# deploy/native/postgrest.service
# systemd unit for PostgREST → /etc/systemd/system/postgrest.service
# Create the service user first:  useradd --system --no-create-home --shell /usr/sbin/nologin postgrest
# Then:  systemctl daemon-reload && systemctl enable --now postgrest
[Unit]
Description=PostgREST for Quire
After=postgresql.service
Wants=postgresql.service

[Service]
User=postgrest
Group=postgrest
ExecStart=/usr/local/bin/postgrest /etc/postgrest/quire.conf
Restart=always
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
