# deploy/native/quire.service
# systemd unit for the Quire app → /etc/systemd/system/quire.service
# Replace the REPLACE_* placeholders (site user, app dir, node bin dir from `which node`
# or a per-user nvm path). The app binds 127.0.0.1:3000; put a reverse proxy in front.
# Then:  systemctl daemon-reload && systemctl enable --now quire
[Unit]
Description=Quire Blog (Next.js)
After=network.target postgrest.service
Wants=postgrest.service

[Service]
Type=simple
User=REPLACE_SITE_USER
Group=REPLACE_SITE_USER
WorkingDirectory=REPLACE_APP_DIR
Environment=PATH=REPLACE_NODE_BIN_DIR:/usr/bin:/bin
Environment=NODE_ENV=production
ExecStart=REPLACE_NODE_BIN_DIR/node node_modules/next/dist/bin/next start -H 127.0.0.1 -p 3000
Restart=always
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
